# scripts/wifi_loader.service
[Unit]
Description=Start MT6620 combo loader

[Service]
Type=oneshot
ExecStart=/usr/sbin/mt6620_loader.sh
RemainAfterExit=true
ExecStop=
