# scripts/wifi_loader.service
Description=Start MT6620 wireless

[Service]
Type=simple
ExecStart=/usr/sbin/mt6620_start_wireless
ExecStop=
WantedBy=network.target
